# nginx/docker/supervisor.d/99-nginx.ini
[program:nginx]
command=/etc/supervisor.d/99-nginx
autorestart=true
