# salt-master/docker/supervisor.d/00-accept-keys.ini
[program:00-accept-keys]
command=/etc/supervisor.d/00-accept-keys
exitcodes=0
startsecs=0
