# salt-master/docker/supervisor.d/00-update-hosts.ini
[program:00-update-hosts]
command=/etc/supervisor.d/00-update-hosts
exitcodes=0
startsecs=0
