# salt-master/docker/supervisor.d/01-salt-master.ini
[program:salt-master]
command=/usr/bin/salt-master --daemon
autorestart=true
