# salt-master/docker/supervisor.d/02-salt-minion.ini
[program:salt-minion]
command=sh -c 'sleep 10 && //usr/bin/salt-minion --log-level debug --daemon'
autorestart=false
