# salt-master/docker/supervisor.d/02-salt-minion.ini
[program:salt-minion]
command=/usr/bin/salt-minion --daemon
autorestart=true
