# elasticsearch/docker/supervisor.d/99-elasticsearch.ini
[program:elasticsearch]
command=gosu elasticsearch /opt/elasticsearch/bin/elasticsearch
autorestart=true
