# elasticsearch/docker/supervisor.d/99-elasticsearch.ini
[program:99-elasticsearch]
command="/etc/supervisor.d/99-elasticsearch"
autorestart=true
